# deploy/titan-core.service
[Unit]
Description=Titan-Core v4.9.6 Sentinel-Vanguard
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/titan-project
ExecStart=/home/ubuntu/titan-project/venv/bin/python3 /home/ubuntu/titan-project/trading_daemon.py

# Redémarrage automatique
Restart=always
RestartSec=30

# Variables d'Environnement (À COMPLÉTER SUR LE SERVEUR)
Environment=PYTHONUNBUFFERED=1
Environment=ENV_MODE=PAPER
Environment=LIVE_AFFIRMATION=False
Environment=ALPACA_API_KEY=VOTRE_CLE
Environment=ALPACA_SECRET_KEY=VOTRE_SECRET
Environment=OPENROUTER_API_KEY=VOTRE_CLE

[Install]
WantedBy=multi-user.target
