# deploy/titan-core.service
[Unit]
Description=Titan-Core v5.6.21 APEX-ULTIMATE
After=network.target

[Service]
# On utilise l'utilisateur 'ubuntu' par défaut sur AWS
User=ubuntu
WorkingDirectory=/home/ubuntu/titan-project
ExecStart=/home/ubuntu/titan-project/venv/bin/python /home/ubuntu/titan-project/trading_daemon.py

# Stratégie de redémarrage : on attend 30s pour laisser respirer les API
Restart=always
RestartSec=30

# --- SECRETS & ENVIRONNEMENT ---
Environment=PYTHONUNBUFFERED=1
Environment=ENV_MODE=PAPER
Environment=LIVE_AFFIRMATION=False
Environment=ALPACA_API_KEY=VOTRE_CLE_ICI
Environment=ALPACA_SECRET_KEY=VOTRE_SECRET_ICI
Environment=FMP_API_KEY=VOTRE_FMP_ICI
Environment=OPENROUTER_API_KEY=VOTRE_OR_ICI
Environment=ALPHA_VANTAGE_KEY=VOTRE_CLE_AV
Environment=TITAN_DASHBOARD_TOKEN=VOTRE_CLE_ICI

# Sécurité système
NoNewPrivileges=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
